# body/bridge/systemd/bridge-watchdog.timer
[Unit]
Description=Run Noema bridge watchdog every 60 seconds

[Timer]
OnBootSec=60s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=bridge-watchdog.service
Persistent=false

[Install]
WantedBy=timers.target
