# body/bridge/systemd/bridge-cycle.service
[Unit]
Description=Noema bridge cycle
Documentation=file:/home/fiste/Noema/bridge/PLAN_bridge_agent.md
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=fiste
WorkingDirectory=/home/fiste/Noema
EnvironmentFile=-/home/fiste/.config/noema/github.env
ExecStart=/usr/bin/python3 /home/fiste/Noema/bridge/scripts/bridge_cycle.py --commit-and-push
StandardOutput=journal
StandardError=journal
SyslogIdentifier=noema-bridge-cycle
TimeoutStartSec=90
Restart=no
Nice=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/home/fiste/Noema

[Install]
WantedBy=multi-user.target
